# templates/boundary.service.erb
[Unit]
Description=<%= scope['boundary::service_name'] %>
After=network.target

[Service]
Type=simple
ExecStart=<%= scope['boundary::install_dir'] %>/boundary > /dev/null 2>&1
ExecStop=<%= scope['boundary::install_dir'] %>/boundary > /dev/null 2>&1
